# util/test_server/.config/acter-mailhog-ci/mailhog.conf
[program:mailhog]
command=/usr/local/bin/mailhog -auth-file=/etc/mailhog-auth.txt
autorestart=true
startretries=10